# build/common-roles/shirasagi.centos7/templates/shirasagi-unicorn.service.j2
[Unit]
Description=Shirasagi Unicorn Server
After=mongod.service

[Service]
EnvironmentFile=-/etc/sysconfig/rvm
User={{ app_user }}
WorkingDirectory={{ install_path }}
ExecStart={{ rvm_home }}/rubies/ruby-{{ rvm_ruby_version }}/bin/bundle exec rake unicorn:start
ExecStop={{ rvm_home }}/rubies/ruby-{{ rvm_ruby_version }}/bin/bundle exec rake unicorn:stop
ExecReload={{ rvm_home }}/rubies/ruby-{{ rvm_ruby_version }}/bin/bundle exec rake unicorn:restart
Type=forking
PIDFile={{ install_path }}/tmp/pids/unicorn.pid

[Install]
WantedBy=multi-user.target
